# repo/backend/docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

; Laravel application server (php artisan serve)
[program:laravel-server]
command=php /var/www/html/artisan serve --host=0.0.0.0 --port=8000
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/www/html/storage/logs/laravel-server.log
stderr_logfile=/var/www/html/storage/logs/laravel-server-error.log

; Laravel queue worker (database driver)
[program:laravel-queue]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600
user=www-data
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
stdout_logfile=/var/www/html/storage/logs/queue-worker.log
stderr_logfile=/var/www/html/storage/logs/queue-worker-error.log

; Cron for Laravel scheduler
[program:cron]
command=crond -f -l 8
user=root
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
